# nat-traversal-notify@.service
[Unit]
Description=NAT Traversal Notification Service

[Service]
Type=oneshot
Environment=SCRIPT_DIR=/home/shuqixiao/Projects/nat-traversal
ExecStart=/bin/bash -c 'source $SCRIPT_DIR/utils.sh && email_notify_on_failure "NAT traversal service (.env.%i) failed: Start limit reached. $(cat $SCRIPT_DIR/.env.%i)"'
